# netmaker/services/netmaker/run
#!/usr/bin/with-contenv sh

echo >&2 "starting netmaker"

touch /app/data/netmaker.db

cd /app

exec ./netmaker 2>&1
