# v2ray/services.d/watcher/run
#!/usr/bin/with-contenv bash
echo >&2 "starting watcher"

sleep 2
exec watchexec -p -w /etc/nginx -- reload-nginx 2>&1
